# systemd/.config/systemd/user/sync-system-font.path
[Unit]
Description=Watch Alacritty config for font changes
After=graphical-session.target

[Path]
PathModified=%h/.config/alacritty/alacritty.toml
Unit=sync-system-font.service

[Install]
WantedBy=default.target